# smokeping/etc/s6-overlay/s6-rc.d/smokeping/run
#!/command/execlineb -P
# smokeping --nodaemon --debug
smokeping --nodaemon
