# smokeping/etc/s6-overlay/s6-rc.d/apache/run
#!/command/execlineb -P
apachectl -D FOREGROUND
